# deploy/modem-manager.service
[Unit]
Description=Modem manager (Huawei MU509 - power, PPP internet, AT control)
After=network.target
StartLimitIntervalSec=0

[Service]
# Run as a normal user, not root - just needs group membership for the
# serial ports (dialout), pppd (dip), and GPIO (gpio). See README "Modem
# daemon setup" for the exact `usermod` commands.
User=root
WorkingDirectory=/home/pisms/sms-dashboard
ExecStart=/home/pisms/sms-dashboard/venv/bin/python3 /home/pisms/sms-dashboard/run_modem_manager.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
